# Apps/806677.com/application/commands/crontab
0 */4 * * * /d/web/yaf/command 806677.com fragment/recacheHotVideo
0 */4 * * * /d/web/yaf/command 806677.com fragment/recacheViewsRankVideo